# com.crystalimages.watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.crystalimages.watcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Developer/CommandLineTools/Library/Frameworks/Python3.framework/Versions/3.9/Resources/Python.app/Contents/MacOS/Python</string>
        <string>/Users/labserver/photo-platform/watch_photos.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/labserver/photo-platform</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/labserver/photo-platform/watcher.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/labserver/photo-platform/watcher.log</string>
</dict>
</plist>
